# supervisord-monitor/py_job_process.conf
[program:test_process]
command=python -u py_job_test.py
directory=/apps/jobs/
stdout_logfile=/apps/jobs/py_job_test.log
stderr_logfile=/apps/jobs/py_job_test_err.log
redirect_stderr=true